# monitor_cluster_health/emr-monitor.service
# EMR Cluster Monitor - Systemd Service
# Install: sudo cp emr-monitor.service /etc/systemd/system/
# Start: sudo systemctl start emr-monitor
# Enable on boot: sudo systemctl enable emr-monitor

[Unit]
Description=EMR Cluster Continuous Monitor
After=network.target

[Service]
Type=simple
User=hadoop
WorkingDirectory=/home/hadoop
ExecStart=/usr/bin/python3 /home/hadoop/monitor_continuous.py --yarn-url http://localhost:8088 --sustained-period 300 --check-interval 60
Restart=always
RestartSec=10
StandardOutput=append:/var/log/emr-monitor.log
StandardError=append:/var/log/emr-monitor.log

# Resource limits
CPUQuota=10%
MemoryLimit=256M

[Install]
WantedBy=multi-user.target
